# infrastructure/vm-stack/test-application-supervisor.ini

[program:test_application]
directory=/home/{{ user }}/hello_world
command=/home/{{ user }}/.env/bin/gunicorn app:app -b localhost:8000
autostart=true
autorestart=true
stderr_logfile=/var/log/test-application/error.log
stdout_logfile=/var/log/test-application/output.log
