# infrastructure/test-application-supervisor.ini
[program:test_application]
directory=/srv/test-application
command=source /srv/test-application/venv/bin/activate && /srv/test-application/venv/bin/gunicorn tiny_app:app
autostart=true
autorestart=true
stderr_logfile=/var/log/test-application/error.log
stdout_logfile=/var/log/test-application/output.log
